# git-backup.service
[Unit]
Description=Git Backup Sync Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/git-backup
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
